# .migration-backup/deploy/xpresspro.service
[Unit]
Description=XpressPro FX API service
After=network.target

[Service]
# Optional environment file (create /etc/xpresspro.env with APP_HOME=/path/to/app)
EnvironmentFile=-/etc/xpresspro.env
# Working directory (falls back to /var/www/Rebranded-xpfx)
WorkingDirectory=${APP_HOME:-/var/www/Rebranded-xpfx}
# Exec the startup script which handles build + run for local/VPS
ExecStart=/usr/bin/env bash -c "${APP_HOME:-/var/www/Rebranded-xpfx}/start.sh"
Restart=always
RestartSec=5
User=www-data
Group=www-data
KillMode=control-group
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target

# Notes:
# 1. Copy this file to /etc/systemd/system/xpresspro.service and create
#    an environment file at /etc/xpresspro.env with APP_HOME pointing at the
#    application directory if you do not want to edit the unit.
# 2. Example /etc/xpresspro.env:
#    APP_HOME=/home/ubuntu/Rebranded-xpfx
# 3. Enable and start:
#    sudo systemctl daemon-reload
#    sudo systemctl enable xpresspro
#    sudo systemctl start xpresspro[Unit]
Description=XpressPro FX API Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Adjust the User and WorkingDirectory for your VPS environment
User=www-data
WorkingDirectory=/var/www/Rebranded-xpfx
EnvironmentFile=-/etc/xpressfx.env
Environment="NODE_ENV=production"
Environment="PORT=3000"

# The start script will install dev deps for the build step when necessary
ExecStart=/usr/bin/env bash /var/www/Rebranded-xpfx/start.sh
Restart=on-failure
RestartSec=5s
StartLimitIntervalSec=300
StartLimitBurst=10

# Security and resource limits
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/www/Rebranded-xpfx /var/log/xpresspro
LimitNOFILE=65536
LimitNPROC=32768
MemoryMax=2G

[Install]
WantedBy=multi-user.target
